# root/etc/services.d/bazarr/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-022}
umask "$UMASK_SET"

exec \
	s6-setuidgid abc python3 /app/bazarr/bazarr.py \
	--no-update --config /config
