# root/etc/services.d/bazarr/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python /app/bazarr/bazarr.py \
--no-update --config /config
